# drive-health/packaging/noctalia-drive-health.timer
[Unit]
Description=Refresh SMART data for Noctalia

[Timer]
OnBootSec=20s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=noctalia-drive-health.service

[Install]
WantedBy=timers.target
